# scripts/launchd/com.obliv.foot.redpay-macstudio-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.obliv.foot.redpay-macstudio-poller
  T-20260711-foot-REDPAY-MACSTUDIO-POLLER — 레드페이 foot 적재 우회로 (맥스튜디오 상주 폴러)

  왜: 레드페이 WAF 가 클라우드/데이터센터 IP(Supabase EF / pg_net) egress 를 403 차단.
      한국 일반 IP 인 맥스튜디오에서만 레드페이 직접 호출이 200 으로 생존(CEO 조종실 소거실험).
      → 이 폴러가 5분 주기로 레드페이 payments.php 직접 조회 → redpay_raw_transactions upsert.

  배포 (supervisor, macstudio):
    1. 시크릿 배치 (평문 커밋 금지, .env* gitignore):
         cat >> ~/.env.redpay-foot <<'EOF'
         SUPABASE_URL="https://rxlomoozakkjesdqjtvd.supabase.co"
         SUPABASE_SERVICE_ROLE_KEY="<foot service_role key>"
         REDPAY_API_KEY="<레드페이 파트너 키>"
         REDPAY_BUSINESS_NO="511-60-00988"
         # ── 화이트리스트 (T-20260711-REDPAY-TERMINAL-REGISTRY-TABLE) ──────────────────────
         #   SSOT = DB 테이블 redpay_terminal_registry(domain=foot,active). 폴러가 자동 조회 → 아래 2줄은
         #   더 이상 필수 아님(생략 시 DB SSOT 파생). 아래를 설정하면 그 값이 DB보다 우선(비상 override 전용).
         #   단말 추가/교체는 env 가 아니라 registry 테이블 seed 1곳만 갱신하면 8곳이 자동 반영(drift 봉인).
         # (선택) 1차 권위 = merchant_id 26(도메인 경계, 17→26 T-20260720-...-WHITELIST-EXPAND). 비워두면 DB registry SSOT 파생.
         # REDPAY_MERCHANT_WHITELIST="1777285001,1777285003,1777285004,1777285005,1777285006,1777285007,1777285008,1777288001,1777288003,1777288004,1777288005,1777288006,1777288008,1777289001,1777289002,1777289003,1777289004,1777289005,1777289006,1777289007,1777289008,1777289009,1777289010,1777289011,1777289012,1777289013"
         # (선택) 보조 TID 26. 비워두면 DB registry SSOT 파생.
         # REDPAY_TID_WHITELIST="1047479255,1047479254,1047479261,1047479268,1047479262,1047479263,1047479264,1047479469,1047479471,1047479472,1047479473,1047479474,1047479475,1047479483,1047479476,1047479477,1047479478,1047479479,1047479480,1047479481,1047479482,1047479153,1047479148,1047479155,1047479158,1047479157"
         EOF
         chmod 600 ~/.env.redpay-foot
    2. 등록:
         mkdir -p ~/logs
         ln -sf ~/GitHub/obliv-foot-crm/scripts/launchd/com.obliv.foot.redpay-macstudio-poller.plist \
           ~/Library/LaunchAgents/com.obliv.foot.redpay-macstudio-poller.plist
         launchctl load -w ~/Library/LaunchAgents/com.obliv.foot.redpay-macstudio-poller.plist
    3. 검증:
         launchctl list | grep redpay-macstudio-poller
         tail -50 ~/logs/redpay_macstudio_poller.out   # '✅ 레드페이 200 OK' + upserted=N 확인
         tail -50 ~/logs/redpay_macstudio_poller.err
-->
<plist version="1.0">
<dict>
    <key>Label</key><string>com.obliv.foot.redpay-macstudio-poller</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/domas/GitHub/obliv-foot-crm &amp;&amp; node scripts/redpay_macstudio_poller.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/domas/GitHub/obliv-foot-crm</string>

    <!-- 5분(300초) 주기 증분 폴 -->
    <key>StartInterval</key><integer>300</integer>
    <key>RunAtLoad</key><true/>

    <!-- SPOF 완화(AC-7): 크래시 시 자동 재기동. 정상 종료(exit 0)는 다음 StartInterval 대기. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key><false/>
    </dict>
    <!-- 재기동 폭주 방지: 실패 후 최소 60초 간격 -->
    <key>ThrottleInterval</key><integer>60</integer>

    <key>StandardOutPath</key><string>/Users/domas/logs/redpay_macstudio_poller.out</string>
    <key>StandardErrorPath</key><string>/Users/domas/logs/redpay_macstudio_poller.err</string>
</dict>
</plist>
